# services/bluetooth-manager.service
[Unit]
Description=Bluetooth Connection Manager
Documentation=Manages Bluetooth device pairing, trusting, and connection monitoring
After=bluetooth.service network.target
Wants=bluetooth.service
Requires=bluetooth.service

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/pi-radio
ExecStart=/home/pi/pi-radio/bluetooth-manager.sh
Restart=always
RestartSec=10
StartLimitInterval=300
StartLimitBurst=5

# Environment
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HOME=/home/pi

# Security settings
NoNewPrivileges=true
ProtectHome=read-only
ProtectSystem=strict
ReadWritePaths=/home/pi/pi-radio /tmp /var/log

# Process settings
KillMode=mixed
KillSignal=SIGTERM
TimeoutStopSec=30

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=bluetooth-manager

[Install]
WantedBy=multi-user.target
