# services/connect-speaker.service
[Unit]
Description=Bluetooth Connection Manager
After=bluetooth.service network.target
Wants=bluetooth.service
Requires=bluetooth.service
StartLimitIntervalSec=60s
StartLimitBurst=2

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/radio-pi
ExecStart=/home/pi/radio-pi/scripts/connect-speaker.sh
Restart=always
#RestartSec=15

# Environment
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=HOME=/home/pi

# Security settings
ReadWritePaths=/home/pi/radio-pi/logs /tmp /var/log

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=connect-speaker

[Install]
WantedBy=multi-user.target
